# systemd/overseer.service
[Unit]
Description=Environment Overseer Daemon

[Service]
ExecStart=overseer

[Install]
WantedBy=multi-user.target
