# systemd/overseer-api.service
[Unit]
Description=Environment Overseer Daemon

[Service]
ExecStart=overseer-rest-api.py

[Install]
WantedBy=multi-user.target
